# system_files/usr/lib/systemd/system/armada-decky-watchdog.service
[Unit]
Description=Restart Decky Loader if it has stopped
# Decky's own bundled PluginLoader (fetched pre-built from upstream at image
# build time, not something armada can patch) has a self-shutdown path that
# fires after repeated CEF disconnects in a short window - normal after a
# session restart (Switch to Desktop, Return to Game Mode, the display
# menu's "Apply & Restart") - and that path calls its own `systemctl stop`
# on itself. Confirmed live: plugin_loader.service's own Restart=always
# never fires afterward (an explicit stop, even self-issued, is never
# auto-restarted by systemd, by design) - Decky just stays gone until
# something starts it back up by hand.
After=plugin_loader.service

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c 'systemctl is-active --quiet plugin_loader.service || systemctl start plugin_loader.service'
