# system_files/usr/lib/systemd/system/armada-first-boot-reboot.service
[Unit]
Description=Reboot once on first boot, after device quirks apply and before initial setup
Documentation=man:systemd.unit(5)
# ConditionFirstBoot checks for a fresh/unset machine-id, which systemd
# itself populates during this same boot - so this naturally fires exactly
# once (the reboot it triggers lands on a boot where machine-id is already
# set, so the condition is false from then on) with no custom marker file.
ConditionFirstBoot=yes
After=armada-device-quirks.service
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reboot

[Install]
WantedBy=multi-user.target
